# hosts/android_server/var/service/ttsd/run
#!/data/data/com.termux/files/usr/bin/bash
exec termux-tts-server --foreground 2>&1
